# gcp-cloud-engineering-portfolio/linux-advanced-troubleshooting-projects/00-systemd-service/demo-api.service
[Unit]
Description=Demo API Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m http.server 8080
WorkingDirectory=/opt/demo-api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

